# deploy/royal-backup.timer
[Unit]
Description=Roda o backup diario da Royal (royal-backup.service)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
